# debian/uucp.socket
[Unit]
Description=UUCP server activation socket

[Socket]
ListenStream=540
Accept=true

[Install]
WantedBy=sockets.target
